# bridge/ai.openclaw.watcher-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.watcher-bridge</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/nacoleon/HQ/02-Projects/SenseCap Watcher/bridge/dist/index.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/nacoleon/HQ/02-Projects/SenseCap Watcher/bridge</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/nacoleon/.openclaw/logs/watcher-bridge.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/nacoleon/.openclaw/logs/watcher-bridge.log</string>
</dict>
</plist>
